# deploy/systemd/uab-digest.timer
[Unit]
Description=ultra-agents-brain daily digest timer
Requires=uab-brain.service

[Timer]
OnCalendar=*-*-* 20:00:00
Persistent=true

[Install]
WantedBy=timers.target
