# ops/com.ultraagents.vault-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ultraagents.vault-sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/caiobellizzi/Documents/Projects/ultra-agents-brain/ops/sync-vault-to-vps.sh</string>
    </array>

    <key>StartInterval</key>
    <integer>300</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/uab-vault-sync.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/uab-vault-sync.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
